# cron/job.cron
*	15	*	*	1,2,3,4,5	/root/public_html/scrapers/cron/scripts/munchery_sf.sh
*	15	*	*	1,2,3,4,5	/root/public_html/scrapers/cron/scripts/munchery_seattle.sh
*	15	*	*	1,2,3,4,5	/root/public_html/scrapers/cron/scripts/sprig.sh

